# roles/i3/files/i3-session.target
 lines (9 sloc) 212 Bytes
[Unit]
Description=i3 session
BindsTo=graphical-session.target